# hostagent/launchd/com.ultraclaude.hostagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.ultraclaude.hostagent</string>
  <!-- Logs go to ~/hostagent.log via the exec redirect inside the zsh -lc
       string, NOT StandardOutPath/StandardErrorPath: launchd does not expand
       ~ in those keys, and adopt-host.sh copies this plist verbatim to hosts
       with arbitrary usernames, then tails ~/hostagent.log for the
       first-heartbeat line. -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>exec >> ~/hostagent.log 2>&amp;1; set -a; source ~/hostagent.env; set +a; exec ~/.bun/bin/bun run ~/hostagent/src/index.ts</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
</dict>
</plist>
